# qml/python/systemd/WiFiKilL3r.timer
[Unit]
Description=WiFiKilL3r Cronjob timer

[Timer]
OnUnitActiveSec=60s

[Install]
WantedBy=timers.target
